# ops/paperclip-vps/hermes-gateway.service
# /etc/systemd/system/hermes-gateway.service
[Unit]
Description=Hermes Agent gateway (worker Nous pour Paperclip)
After=network-online.target paperclip.service
Wants=network-online.target

[Service]
Type=simple
User=paperclip
WorkingDirectory=/home/paperclip
ExecStart=/home/paperclip/.local/bin/hermes gateway start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
